# dev/prodexui.service
[Unit]
Description=prodexui

[Service]
Type=simple
ExecStart=%h/etc/prodex/prodex -mode ui
Restart=always

[Install]
WantedBy=multi-user.target
