# dev/prodex.timer
[Unit]
Description=prodex
Requires=prodex.service

[Timer]
Unit=prodex.service
OnCalendar=*-*-* *:00:00

[Install]
WantedBy=timers.target
